# lib/supervisord.conf
[program:httpd]
priority = 100
command = /bin/bash -c "/bin/sleep 2 && /usr/sbin/httpd -D FOREGROUND"
redirect_stderr = true
stdout_events_enabled = true
startsecs = 5
autorestart = false

[program:php-fpm]
priority = 200
command = /bin/bash -c "/bin/sleep 2 && /usr/sbin/php-fpm -c /etc/php.ini"
stdout_events_enabled = true
startsecs = 5
autorestart = false

[unix_http_server]
username = superuser
password = superuser

[supervisorctl]
username = superuser
password = superuser

[supervisord]
user = root
